# devel/cockpit-container/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=9091
